# scripts/openvpn-uds-monitor.service
[Unit]
Description=Optimized OpenVPN UDS Traffic Monitor Service
After=network.target openvpn@server-cert.service openvpn@server-login.service
Wants=openvpn@server-cert.service
PartOf=openvpn@server-cert.service

[Service]
Type=simple
User=root
Group=openvpn
WorkingDirectory=/root/openvpn
ExecStart=/usr/bin/python3 /root/openvpn/scripts/uds_monitor_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/openvpn /root/openvpn/openvpn_data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

LimitNOFILE=65536
MemoryMax=50M
CPUQuota=5%

EnvironmentFile=-/root/openvpn/.env
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
