# scripts/openvpn-uds-monitor.service
[Unit]
Description=OpenVPN UDS Traffic Monitor Service
After=network.target openvpn-server@server-cert.service openvpn-server@server-login.service
Wants=openvpn-server@server-cert.service
PartOf=openvpn-server@server-cert.service

[Service]
Type=simple
User=root
Group=openvpn
WorkingDirectory=/etc/owpanel
ExecStart=/usr/bin/python3 /etc/owpanel/scripts/uds_monitor_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/openvpn /etc/owpanel/openvpn_data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
MemoryMax=100M
CPUQuota=10%

# Environment configuration
EnvironmentFile=/etc/owpanel/.env
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 