# api/openvpn-api.service
[Unit]
Description=OpenVPN Manager API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/openvpn
Environment=PATH=/root/openvpn/venv/bin
Environment=PYTHONPATH=/root/openvpn
EnvironmentFile=/etc/openvpn-manager/.env
EnvironmentFile=/root/openvpn/environment.env

ExecStart=/root/openvpn/venv/bin/python api/app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/openvpn/openvpn_data /var/log/openvpn

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openvpn-api

[Install]
WantedBy=multi-user.target