# api/openvpn-api.service
[Unit]
Description=OpenVPN Manager API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/seyed/Cursor Project/openvpn
Environment=PYTHONPATH=/home/seyed/Cursor Project/openvpn
Environment=OPENVPN_API_KEY=your-secure-api-key-here
ExecStart=/usr/bin/python3 /home/seyed/Cursor Project/openvpn/api/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target