# api/openvpn-api.service
[Unit]
Description=OpenVPN Manager API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/etc/owpanel
Environment=PATH=/etc/owpanel/venv/bin
Environment=PYTHONPATH=/etc/owpanel
EnvironmentFile=/etc/owpanel/.env

ExecStart=/etc/owpanel/venv/bin/python api/app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/owpanel/openvpn_data /var/log/openvpn

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openvpn-api

[Install]
WantedBy=multi-user.target